# systemd/system-upgrade-daemon.service
[Unit]
Description=System Upgrade Daemon
After=network.target

[Service]
Type=forking
ExecStart=/usr/bin/system-upgrade-daemon

[Install]
WantedBy=multi-user.target
